# magickey.service
[Unit]
Description=Magickey -- remapping keys

[Service]
Type=simple
ExecStart=magickey

